# praktikum-modul-1-d04/task-3/crontabs
* * * * * /usr/operating-system/praktikum-modul-1-d04/task-3/minute_log.sh
@hourly /usr/operating-system/praktikum-modul-1-d04/task-3/aggregate_minutes_to_hourly_log.sh
59 23 * * * /usr/operating-system/praktikum-modul-1-d04/task-3/bakcup_metrics.sh
